# deploy/cpp-tutor-backend.service
[Unit]
Description=cpp-tutor FastAPI backend
After=network.target docker.service
Requires=docker.service

[Service]
User=azureuser
WorkingDirectory=/home/azureuser/cpp-tutor/backend
Environment=CPP_TUTOR_CORS_ORIGINS=https://cpp-tutor.centralindia.cloudapp.azure.com
ExecStart=/home/azureuser/cpp-tutor/backend/.venv/bin/uvicorn app.api:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
